# voice-cloning-api.service
[Unit]
Description=Voice Cloning API Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/project
Environment="PATH=/path/to/your/venv/bin"
ExecStart=/path/to/your/venv/bin/gunicorn -c gunicorn_config.py api:app
Restart=always

[Install]
WantedBy=multi-user.target 